# logdaemon/units/logdaemon.service.in
[Unit]
Description=Logging Service Tutorial

[Service]
ExecStart=@rootbindir@/logdaemon

[Install]
WantedBy=multi-user.target
