# sadaemon/units/sadaemon.socket.in
[Unit]
Description=Socket-Activated Service Tutorial Socket

[Socket]
ListenStream=23

[Install]
WantedBy=sockets.target

