# sadaemon/units/sadaemon.service.in
[Unit]
Description=Socket-Activated Service Tutorial

[Service]
ExecStart=@rootbindir@/sadaemon

